# infra/systemd/regex-app.service.example
[Unit]
Description=Word Regex Search gunicorn service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/word-regex-search/apps/regex/src
Environment="PATH=/var/www/word-regex-search/apps/regex/.venv/bin"
Environment="WORD_LIST_PATH=/usr/share/dict/words"
Environment="RESULT_LIMIT=200"
Environment="REGEX_TIMEOUT_SECONDS=0.05"
Environment="MAX_PATTERN_LENGTH=256"
ExecStart=/var/www/word-regex-search/apps/regex/.venv/bin/gunicorn --workers 3 --bind unix:/run/word-regex-search.sock --timeout 30 --max-requests 1000 --max-requests-jitter 100 wsgi:app
Restart=always
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/var/www/word-regex-search/apps/regex
ReadWritePaths=/run

[Install]
WantedBy=multi-user.target
